# Linux/TDSQL/dump-database.service
[Unit]
Description=database

[Service]
ExecStart=/home/brice/Documents/TDSQL/script.sh

[Install]
WantedBy=multi-user.target
