# Linux/TDSQL/dump-database.timer
[Unit]
Description=sauvegarde 10min

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
