# Linux/TDSQL/dump-database.service

````shell
[Unit]
Description=database

[Service]
ExecStart=/usr/bin/script.sh

[Install]
WantedBy=multi-userr.target
````