# Linux/TDSQL/dump-database.service
[Unit]
Description=database

[Service]
ExecStart=/usr/bin/script.sh

[Install]
WantedBy=multi-user.target
